# templates/default/consul-init.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          consul
# Required-Start:    $network $remote_fs $syslog
# Required-Stop:     $network $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Consul Service Discovery Platform
# Description:       Consul is a tool for discovering and configuring services
#                    in your infrastructure. It provides several key features:
#                       * Service Discovery
#                       * Health Checking
#                       * Key/Valuye Store
#                       * Multi Datacenter
### END INIT INFO

export GOMAXPROCS=${GOMAXPROCS}

DAEMON="<%= @consul_binary %>"
DAEMON_ARGS="agent -config-dir <%= @config_dir %> 2> $LOGFILE"
NAME="consul"

PIDFILE="/var/run/$NAME.pid"
LOGFILE="/var/log/$NAME.log"

do_start() {
  start-stop-daemon \
    --start \
    --background \
    --exec $DAEMON \
    --make-pidfile \
    --pidfile $PIDFILE \
    -- $DAEMON_ARGS
}

do_stop() {
  start-stop-daemon \
    --stop \
    --pidfile $PIDFILE
}

do_reload() {
  start-stop-daemon \
    --exec $DAEMON \
    --pidfile $PIDFILE \
    -s 1
}

do_status() {
  start-stop-daemon \
    --status
    --pidfile $PIDFILE
}

case "$1" in
  start)
    echo "Starting $NAME"
    do_start
  ;;
  stop)
    echo "Stopping $NAME"
    do_stop
  ;;
  restart)
    echo "Restarting $NAME"
    do_stop
    do_start
  ;;
  reload)
    echo "Reloading $NAME"
    do_reload
  ;;
  status)
    echo "Checking $NAME status"
    do_status
  ;;
  *)
    echo "Usage: $0 {start|stop|restart|reload|status}"
    exit 1
  ;;
esac

exit 0
